# go-bot/deploy/xray-telegram-bot.service
[Unit]
Description=Xray Telegram Bot (Go)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/usr/local/xray-cli/go-bot
EnvironmentFile=-/etc/default/xraycli-api
EnvironmentFile=-/etc/default/xraycli-telegram-bot
ExecStart=/usr/local/xray-cli/go-bot/bin/xray-telegram-bot
Restart=always
RestartSec=3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/usr/local/xray-cli
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
